# DEBIAN/postinst
#!/bin/bash

chmod +x /opt/evalgo.org/evhetzner-service/bin/evhetzner
systemctl enable supervisor
systemctl start supervisor
supervisorctl reload
